# scripts/com.user.hugging-papers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.hugging-papers</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sweden/Desktop/Pandora/Projects/hugging_papers/scripts/run_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/sweden/Desktop/Pandora/Projects/hugging_papers</string>

    <key>StandardOutPath</key>
    <string>/Users/sweden/Desktop/Pandora/Projects/hugging_papers/logs/launchd.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/sweden/Desktop/Pandora/Projects/hugging_papers/logs/launchd.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
